# ansible/roles/app/templates/mywebapp.socket.j2
[Unit]
Description=mywebapp socket

[Socket]
ListenStream=5200
Accept=no

[Install]
WantedBy=sockets.target